# loffeblogg/systemd/loffeblogg.timer
[Unit]
Description=Loffeblogg auto-update timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
