# loffeblogg/systemd/loffeblogg.service
[Unit]
Description=Loffeblogg auto-update and deploy
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%WORKDIR%
ExecStart=%WORKDIR%/scripts/auto-update.sh
Environment="PATH=%HOME%/bin:%HOME%/.local/bin:%HOME%/.nvm/versions/node/v22.17.1/bin:/usr/local/bin:/usr/bin:/bin"
Environment="SSH_AUTH_SOCK=/run/user/1000/keyring/ssh"

[Install]
WantedBy=default.target
