# babysovelogg/deploy/systemd/babysovelogg@.service
[Unit]
Description=babysovelogg instance for %i
After=network.target
AssertPathExists=/srv/babysovelogg/code/build/index.js
AssertPathExists=/srv/babysovelogg/families/%i

[Service]
Type=simple
User=babysovelogg
Group=babysovelogg
UMask=0077
WorkingDirectory=/srv/babysovelogg/families/%i
# .env holds NODE_ENV/ORIGIN/VAPID_*/SOCKET_PATH/DB_PATH — managed by ansible.
# (Keep DB_PATH/SOCKET_PATH there, not in this unit: EnvironmentFile overrides
# Environment=, so duplicating here would let a stray .env value win.)
EnvironmentFile=/srv/babysovelogg/families/%i/.env
ExecStartPre=/usr/bin/rm -f /srv/babysovelogg/families/%i/app.sock
ExecStart=/usr/local/bin/bun /srv/babysovelogg/code/build/index.js
ExecStartPost=/bin/sh -c 'until [ -S /srv/babysovelogg/families/%i/app.sock ]; do sleep 0.05; done; chmod 0660 /srv/babysovelogg/families/%i/app.sock'
Restart=on-failure
RestartSec=5

# Sandboxing — meetings.service shape.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/srv/babysovelogg/families/%i
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
