# babysovelogg/deploy/systemd/babysovelogg@.service
[Unit]
Description=babysovelogg instance for %i
After=network.target
AssertPathExists=/srv/babysovelogg/code/build/index.js
AssertPathExists=/srv/babysovelogg/families/%i

[Service]
Type=simple
User=babysovelogg
Group=babysovelogg
UMask=0077
WorkingDirectory=/srv/babysovelogg/families/%i
# .env holds NODE_ENV/ORIGIN/HOST/PORT/VAPID_*/DB_PATH — managed by ansible.
EnvironmentFile=/srv/babysovelogg/families/%i/.env
ExecStart=/usr/local/bin/bun /srv/babysovelogg/code/build/index.js
Restart=on-failure
RestartSec=5
# Bun's adapter-node graceful shutdown never lets the process exit
# (Bun bug: httpServer.close() callback doesn't fire after
# closeAllConnections). So we accept SIGKILL after 5s — the
# hooks.server.ts cleanup runs on SIGTERM when hooks are loaded, which
# closes the DB cleanly before the kill. WAL-mode SQLite handles the
# rest. Lower TimeoutStopSec if 5s feels too long; raise to allow
# slow in-flight requests to complete.
TimeoutStopSec=5

# Sandboxing — meetings.service shape.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/srv/babysovelogg/families/%i
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
